# systemd/unblock_pihole.service
[Unit]
Description=Pi-hole Control Panel
After=network.target pihole.service

[Service]
Type=simple
User=unblock-pihole
Group=unblock-pihole
WorkingDirectory=/opt/unblock_pihole
EnvironmentFile=/opt/unblock_pihole/.env
ExecStart=/opt/unblock_pihole/.venv/bin/python -m unblock_pihole

# Restart on failure
Restart=on-failure
RestartSec=10

# This ensures signals are passed through to Python properly
KillMode=mixed
KillSignal=SIGTERM

TimeoutStopSec=15s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/opt/unblock_pihole

[Install]
WantedBy=multi-user.target